# systemd/arpCheck.service
[Unit]
Description=ARP checks
Wants=arpCheck.timer

[Service]
Type=simple
ExecStart=~/scrypt.sh

[Install]
WantedBy=multi-user.target
